# crontab
* * * * * curl https://saltwater-recon-r6ebf.ondigitalocean.app/api/
* * * * * echo 'Hello cron!''

# An empty line is required at the end of this file for a valid cron file.
